# src/Stubs/docker/prod/supervisord/supervisord-30-queue.conf
[program:queue]
command=/usr/local/bin/php artisan mercury:boot-horizon
directory=/srv/www
user=aurora
group=aurora
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
